# supervisor/logstash.conf
[program:logstash]
directory=/
command=/opt/logstash/bin/logstash agent -f /etc/logstash/conf.d
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/logstash.out.log
stderr_logfile=/var/log/supervisor/logstash.err.log
